# utils/systemd/monerod.service
[Unit]
Description=scala Full Node
After=network.target

[Service]
User=scala
Group=scala
WorkingDirectory=~
StateDirectory=scala
LogsDirectory=scala

Type=simple
ExecStart=/usr/bin/scalad --config-file /etc/scalad.conf --non-interactive
StandardOutput=null
StandardError=null

Restart=always

[Install]
WantedBy=multi-user.target
